# deploy/systemd/switchboard.target
[Unit]
Description=Switchboard VoIP Platform
Documentation=https://github.com/sebas/switchboard

# This target groups all switchboard services together.
# Use: systemctl start switchboard.target
#      systemctl stop switchboard.target
#      systemctl restart switchboard.target

# Pull in all switchboard services when this target is started.
# Order matters: rtpmanager must start before signaling.
Wants=switchboard-rtpmanager.service
Wants=switchboard-signaling.service
Wants=switchboard-ui.service

# Ensure proper startup order
After=network.target
After=switchboard-rtpmanager.service
After=switchboard-signaling.service

[Install]
WantedBy=multi-user.target
