# deploy/systemd/switchboard-signaling.service
[Unit]
Description=Switchboard SIP Signaling Service
Documentation=https://github.com/sebas/switchboard
After=network.target switchboard-rtpmanager.service
Requires=switchboard-rtpmanager.service
PartOf=switchboard.target

[Service]
Type=simple
User=switchboard
Group=switchboard
WorkingDirectory=/opt/switchboard

ExecStart=/bin/bash -c '/opt/switchboard/switchboard-signaling \
    --port=${SIP_PORT} \
    --bind=${SIP_BIND} \
    ${SIP_ADVERTISE:+--advertise=${SIP_ADVERTISE}} \
    --rtpmanager=${RTPMANAGER_ADDRS} \
    --dialplan=${DIALPLAN_PATH} \
    --loglevel=${LOG_LEVEL}'

EnvironmentFile=/etc/switchboard/signaling.env

Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/switchboard

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=switchboard-signaling

[Install]
WantedBy=switchboard.target
WantedBy=multi-user.target
