# deploy/systemd/switchboard-rtpmanager.service
[Unit]
Description=Switchboard RTP Manager Service
Documentation=https://github.com/sebas/switchboard
After=network.target
PartOf=switchboard.target

[Service]
Type=simple
User=switchboard
Group=switchboard
WorkingDirectory=/opt/switchboard

ExecStart=/opt/switchboard/switchboard-rtpmanager \
    --grpc-port=${GRPC_PORT} \
    --grpc-bind=${GRPC_BIND} \
    --advertise=${RTP_ADVERTISE} \
    --rtp-port-min=${RTP_PORT_MIN} \
    --rtp-port-max=${RTP_PORT_MAX} \
    --audio-path=${AUDIO_PATH}

EnvironmentFile=/etc/switchboard/rtpmanager.env

Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/switchboard/audio

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=switchboard-rtpmanager

[Install]
WantedBy=switchboard.target
WantedBy=multi-user.target
